# deploy/oracle/birdcam-api.service
[Unit]
Description=BirdCam FastAPI service
After=network.target mysql.service

[Service]
Type=simple
User=birdcam
WorkingDirectory=/opt/birdcam/backend
EnvironmentFile=/opt/birdcam/backend/.env
ExecStart=/opt/birdcam/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
